# weatherlink_restart.timer
[Unit]
Description=Execute task every xx

[Timer]
OnCalendar=*-*-* 00:00:30

Unit=weatherlink_restart.service

[Install]
WantedBy=timers.target 
